# docker/supervisord-mongodb.conf
[program:mongodb]
command=/metis/docker/start-mongodb.sh
numprocs=1
autostart=true
autorestart=true
priority=1
