# docker/supervisord-metis.conf
[program:metis]
command=/metis/docker/start-metis.sh
numprocs=1
autostart=true
autorestart=true
priority=200
