# docker/supervisord-import-seed-data.conf
[program:import-seed-data]
command=/metis/docker/import-seed-data.sh
numprocs=1
autostart=true
autorestart=false
priority=100
